# root/etc/udev/rules.d/90-setupcon.rules
ACTION=="add", SUBSYSTEM=="vtconsole", KERNEL=="vtcon*", RUN+="/bin/setupcon"
